# roles/smartcity/files/mnt-imstorage.automount
[Unit]
Description=Automount Pronghorn Association Storage

[Automount]
Where=/mnt/storage

[Install]
WantedBy=multi-user.target
