# roles/smartcity/files/mnt-imstorage.automount
[Unit]
Description=Automount RCDC

[Automount]
Where=/mnt/imstorage

[Install]
WantedBy=multi-user.target
